# dist/supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
user=mysql
command=/usr/sbin/mysqld
autorestart=true

[program:cowrie]
user=cowrie
directory=/home/cowrie/cowrie
command=bash -c 'export PYTHONPATH=/home/cowrie/cowrie && exec /usr/local/bin/twistd --nodaemon -y cowrie.tac --pidfile var/run/cowrie.pid cowrie'
autorestart=true

[program:ewsposter]
command=bash -c 'sleep 10 && exec /usr/bin/python /opt/ewsposter/ews.py -c /data/ews/conf/ -m kippo -l 60'
autorestart=true
